# drupal8/crontab
*/1 * * * * root wget http://localhost/cron/heartbeat.json
*/15 * * * * root wget http://localhost/api/every_15_minute.json?mode=j5p4v5s5